# cmake/package/debian/postrm
#!/bin/sh

update-rc.d -f $PKG_NAME remove
