# cmake/package/debian/postrm.in
#!/bin/sh

update-rc.d -f @APP_NAME@ remove
