# pkg/gentoo/net-analyzer/rackspace-monitoring-agent/files/init
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

PIDFILE=/var/run/$PKG_NAME.pid

depend() {
    use net
    use logger dns
}

start() {
    ebegin "Starting $PKG_NAME"
    start-stop-daemon --background --start --exec /usr/bin/$PKG_NAME \
      --make-pidfile --pidfile $PIDFILE
    eend $?
}

stop() {
    ebegin "Stopping $PKG_NAME"
    start-stop-daemon --stop --exec /usr/bin/$PKG_NAME \
      --pidfile $PIDFILE
    eend $?
}

restart() {
    stop
    start
}
