# cmake/package/rpm/postinst.in

# This adds the proper /etc/rc*.d links for the script
/sbin/chkconfig --add $PKG_NAME

mkdir -p /var/lib/rackspace-monitoring-agent

# Restart agent on upgrade
if [ "$1" = "2" ] ; then
    /sbin/service $PKG_NAME stop  >/dev/null 2>&1 || :
    /sbin/service $PKG_NAME start >/dev/null 2>&1 || :
fi

