# pkg/monitoring/debian/package.upstart
# rackspace-monitoring-agent

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [!2345]

respawn
exec rackspace-monitoring-agent --production >> /var/log/rackspace-monitoring-agent.log 2>&1
