# pkg/monitoring/debian/postrm
#!/bin/sh

update-rc.d rackspace-monitoring-agent remove
