# pkg/monitoring/debian/postinst
#!/bin/sh
update-rc.d rackspace-monitoring-agent defaults 91 09
