# pkg/monitoring/debian/postrm
#!/bin/sh

update-rc.d -f rackspace-monitoring-agent remove
